# gui/setup/my_services/openbox.conf
[program:openbox]
command=/usr/bin/openbox-session
environment=HOME="/nobody",USER="nobody",DISPLAY=":0"
