# gui/setup/my_services/xrdp-sesman.conf
[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman
